# deploy/pixel-office.service
[Unit]
Description=Pixel Office — AI Agents Dashboard
After=network.target

[Service]
Type=simple
# Replace with your deployment user (e.g. the account that owns the clone)
User=YOUR_LINUX_USER
WorkingDirectory=/home/YOUR_LINUX_USER/agents-workspace/pixel-office

# ⚠️  Do NOT use /usr/bin/node — it may be old or missing.
# Find your node path: which node
# nvm example: /home/YOUR_LINUX_USER/.nvm/versions/node/v22.x.x/bin/node
ExecStart=/bin/bash -c 'exec $(which node) node_modules/.bin/next start -p 3001'

Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
Environment=PORT=3001
# Copy .env.example → .env.local, fill in values, then point here:
EnvironmentFile=/home/YOUR_LINUX_USER/agents-workspace/pixel-office/.env.local

[Install]
WantedBy=multi-user.target
