# deploy/pixel-office.service
[Unit]
Description=Pixel Office — AI Agents Dashboard
After=network.target

[Service]
Type=simple
User=openclaw
WorkingDirectory=/home/openclaw/agents-workspace/pixel-office
ExecStart=/usr/bin/node node_modules/.bin/next start -p 3001
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
Environment=PORT=3001
# Add your env vars below:
# Environment=OPENCLAW_GATEWAY_URL=http://localhost:18789
# Environment=OPENCLAW_TOKEN=your_token_here
# Environment=DATABASE_URL=/home/openclaw/agents-workspace/pixel-office/data/office.db

[Install]
WantedBy=multi-user.target
